# Web.Upnp.Control/upnp-controller-dashboard.service
[Unit]
Description=UPNP Controller Dashboard Service

[Service]
Type=notify
ExecStart=/usr/sbin/testapp

[Install]
WantedBy=multi-user.target